# linux-image-scripts/centos/rc.local
if [ -x '/root/foreman_finish_deploy.sh' ]; then
  /root/foreman_finish_deploy.sh &
fi

